# data/indicator-transfer.service.in
[Unit]
Description=Indicator Transfer Backend
PartOf=graphical-session.target
After=indicators-pre.target

[Service]
ExecStart=@pkglibexecdir@/indicator-transfer-service
